# multimedia/rygel/files/etc/hotplug.d/iface/30-rygel
#!/bin/sh

case "$ACTION" in
	ifup)
		# Experience proved that rygel needed to be restarted if an
		# interface config was updated so that's what we do here
		/etc/init.d/rygel restart
	;;
esac
