# utils/bluez/files/bluez-utils.init
#!/bin/sh /etc/rc.common
# Copyright (c) 2013 Qualcomm Atheros, Inc.

START=62


start() {
        service_start /usr/sbin/bluetoothd
}

stop() {
        service_stop /usr/sbin/bluetoothd
}
